# .config/systemd/user/gpg-agent.service.d/override.conf
[Service]
Environment="GNUPGHOME=%h/.local/share/gnupg"
